# parser-peg/files/node.service
[Unit]
Description=node
After=network-online.target
Wants=network-offline.target
